# deploy/moadir-deploy.service
# Serviço oneshot que checa/aplica deploys do moadir (chamado pelo timer).
# Instalar em /etc/systemd/system/moadir-deploy.service
# Logs: journalctl -u moadir-deploy -f
[Unit]
Description=moadir - deploy pull (checa GitHub e aplica se houver commit novo)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# HOME é necessário: o git lê o safe.directory de ~/.gitconfig e o pnpm usa
# ~/.local/share/pnpm como store no build. Sem isso o systemd roda sem HOME e
# o deploy falha (dubious ownership / pnpm).
Environment=HOME=/root
# flock evita dois deploys sobrepostos (se um build ainda estiver rodando).
ExecStart=/usr/bin/flock -n /run/moadir-deploy.lock /opt/moadir/deploy/pull-deploy.sh
